# deploy/ansible/templates/istota-browser-watchdog.service.j2
[Unit]
Description={{ istota_namespace | capitalize }} browser container watchdog (health check + debounced auto-heal)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/{{ istota_namespace }}-browser-watchdog check
# The watchdog drives docker compose; it runs as root (no User= override).
SyslogIdentifier={{ istota_namespace }}-browser-watchdog
